# weather.timer
[Unit]
Description=Run weather weekly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

