# weather.service
[Unit]
Description=Weather Daemon

[Service]
ExecStart=/usr/sbin/weather

[Install]
WantedBy=graphical.target
